# deploy/aws-ec2/systemd/telehitch-airflow-scheduler.service
[Unit]
Description=Telehitch Apache Airflow scheduler
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=__AIRFLOW_USER__
Group=__AIRFLOW_GROUP__
WorkingDirectory=__REPO_ROOT__
EnvironmentFile=__RUNTIME_ENV__
EnvironmentFile=__APPLICATION_ENV__
ExecStart=__VENV_ROOT__/bin/airflow scheduler
Restart=always
RestartSec=10
TimeoutStopSec=120
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
