# deploy/aws-ec2/systemd/telehitch-airflow-backup.timer
[Unit]
Description=Daily Telehitch Airflow metadata backup

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15m
Unit=telehitch-airflow-backup.service

[Install]
WantedBy=timers.target
